# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

[program:neo4j]
command=/usr/bin/neo4j console
autorestart=true
stopsignal=INT

[program:qdrant]
command=/usr/local/bin/qdrant
autorestart=true
stopsignal=INT

[program:memos]
command=bash -c 'for i in {1..30}; do python3 -c "import socket; s=socket.socket(); s.settimeout(2); s.connect((\"localhost\",6333)); s.close()" 2>/dev/null && break; sleep 1; done; exec uvicorn memos.api.server_api:app --host 0.0.0.0 --port 8100'
autorestart=true
stopsignal=INT
startsecs=10
environment=HF_ENDPOINT="https://hf-mirror.com"
